# debug/php/5/supervisord.conf
[program:php5-fpm]
command=/usr/sbin/php5-fpm --nodaemonize --fpm-config /etc/php5/fpm/php-fpm.conf
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/tmp/php-fpm.log
